# kube-node/deb/kube-node/postinstall
#!/bin/bash

# enable service
systemctl enable containerd \
    kubelet \
    kube-proxy \
    cni-bridge-setup