# kube-controller/deb/kube-controller/postinstall
#!/bin/bash

# enable service
systemctl enable kube-apiserver kube-controller-manager kube-scheduler kube-config