# kube-node/deb/kube-node/debian/etc/systemd/system/kubelet.service
[Unit]
Description=Kubernetes Kubelet
Documentation=https://github.com/kubernetes/kubernetes
After=containerd.service
Requires=containerd.service

[Service]
ExecStart=/bin/bash /etc/kubelet/start-kubelet.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target